# files/usr/lib/systemd/user/flatpak-fetch-updates.timer
[Unit]
Description=Fetch flatpak updates

[Timer]
OnBootSec=1h
OnUnitInactiveSec=8h
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
