# files/usr/lib/systemd/user/flatpak-fetch-updates.service
[Unit]
Description=Fetch flatpak updates

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update --assumeyes --noninteractive
